# file_root/mongodb-base/mongodb-upstart.conf
limit nofile 20000 20000

kill timeout 300

start on runlevel [2345]
stop on runlevel [06]

exec /usr/bin/mongod --config /etc/mongodb.conf

